# nabd/nabd.socket
[Socket]
ListenStream=0.0.0.0:10543

[Install]
WantedBy = sockets.target
